# supervisord.conf
[program:twumblrd]
command = python /home/dotcloud/current/twumblr/scripts/worker.py
stderr_logfile = /var/log/supervisor/awesome.log