# supervisord.conf
[program:twumblrd]
command = python /home/dotcloud/current/twumblr/scripts/worker.py